# src/rotate.cron
01 00 * * * /home/zhao/www/receiver/receiver/src/logrotate
